# client/daily-restart.service
[Unit]
Description=Daily Restart Service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/systemctl --message="Daily Reboot" reboot
