# client/matrix-client.conf
[program:matrix-client]
command=bash -c "sleep 5 && exec /home/pi/client --rows 32 --cols 64 --slowdown 4"
directory=/home/pi/
autostart=true
autorestart=true
